# packaging/deb/postrm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl disable --now openclaw-gateway.service || true
  systemctl daemon-reload || true
fi

exit 0
